# deployment/do_it.conf
# supervisor conf
[program:do_it]
command=/home/do_it/bin/gunicorn_start
user=do_it
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/do_it/logs/gunicorn-error.log
